# images/server/lab-entrypoint.service
[Unit]
Description=Lab Entrypoint Service
After=network.target

[Service]
Type=simple
Environment=LOGIN_USER=admin
Environment=LOGIN_PASSWORD=admin
Environment=DB_USER=normal_user
Environment=DB_PASSWORD=normalpass
Environment=ATTACKER_USER=normal_user
Environment=ATTACKER_PASSWORD=normalpass
Environment=DEF_ROLE=superhero
EnvironmentFile=-/etc/default/lab-env
ExecStart=/entrypoint.sh
Restart=always

[Install]
WantedBy=multi-user.target
